# ansible/roles/services/files/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
stdout_logfile=/var/log/supervisor/httpd.log
stderr_logfile=/var/log/supervisor/httpd.err
autostart=true
autorestart=true
stopsignal=QUIT
