# ansible/roles/services/files/justus.ini
[program:justus]
command=/usr/bin/npm start
directory=/opt/sources
autostart=true
autorestart=true
environment=NODE_ENV=dev
stdout_logfile=/var/log/supervisor/justus.log
stderr_logfile=/var/log/supervisor/justus.err
user=vagrant
