# ansible/roles/services/files/postgres.ini
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/pgsql/data -p 5432
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.err
user=postgres
autostart=true
autorestart=true
stopsignal=QUIT
