# ansible/roles/services/files/redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
directory=/opt/sources
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err
user=root
autostart=true
autorestart=true
stopsignal=QUIT
